# tasks/crontab
* * * * * /usr/local/bin/reset_monthly_usage.sh
